# custom_novel/custom_ai_novel.ini
[program:custom_ai_novel]

command=uwsgi /home/data/custom_ai_novel/uwsgi.ini